# marquee-update.timer
[Unit]
Description=Nightly Cubs Marquee update check (4 AM)

[Timer]
# Overnight-only window so an update reboot never interrupts a game
OnCalendar=*-*-* 04:00:00
# Spread multiple marquees out a bit
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
