# wifi-manager.service
[Unit]
Description=WiFi Manager with AP Fallback for Cubs Scoreboard
After=network.target
Before=cubs-scoreboard.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi
ExecStart=/home/pi/wifi_manager.sh

# Restart configuration
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Give network time to initialize
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target