# wifi-manager.service
[Unit]
Description=WiFi Manager with AP Fallback for Cubs Scoreboard
After=network-pre.target avahi-daemon.service
Wants=network-pre.target avahi-daemon.service
Requires=avahi-daemon.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi
ExecStart=/home/pi/wifi_manager.sh

# Restart configuration
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Give network time to initialize - increased timeout
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
